# root/etc/services.d/rrshareweb/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /rrshare/rrshareweb &> /dev/null &